# deploy/supervisor.conf
; =============================================================
; Supervisor Configuration — bahleel.io / smarthom
; File ini diletakkan di: /etc/supervisor/conf.d/smarthom.conf
;
; Setelah menyalin file ini:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start all
; =============================================================

; -------------------------------------------------------
; [1] MQTT Listener
;     Proses long-running yang subscribe ke MQTT broker,
;     update widget value, dan catat log ke database.
; -------------------------------------------------------
[program:laravel-mqtt]
process_name=%(program_name)s
command=php /var/www/smarthom/artisan mqtt:listen
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=10
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-mqtt.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=HOME="/var/www",USER="www-data"

; -------------------------------------------------------
; [2] Laravel Scheduler (daemon)
;     Menggantikan cron `* * * * * php artisan schedule:run`.
;     schedule:work secara internal loop tiap menit.
; -------------------------------------------------------
[program:laravel-scheduler]
process_name=%(program_name)s
command=php /var/www/smarthom/artisan schedule:work
autostart=true
autorestart=true
startsecs=3
startretries=5
stopwaitsecs=5
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
environment=HOME="/var/www",USER="www-data"

; -------------------------------------------------------
; [3] Laravel Queue Worker
;     Memroses jobs di queue (database driver).
;     --max-time=3600 agar restart otomatis tiap jam
;     (menghindari memory leak pada proses panjang).
; -------------------------------------------------------
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/smarthom/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --queue=default
autostart=true
autorestart=true
numprocs=1
startsecs=5
startretries=5
stopwaitsecs=60
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/laravel-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=HOME="/var/www",USER="www-data"

; -------------------------------------------------------
; Group — restart semua sekaligus dengan:
;   sudo supervisorctl restart smarthom:*
; -------------------------------------------------------
[group:smarthom]
programs=laravel-mqtt,laravel-scheduler,laravel-queue
priority=999
