# deploy/supervisor.conf
; =============================================================
; Supervisor Configuration — bahleel.io / smarthom
; Server path: /www/wwwroot/nh.mdpower.io/smarthom
;
; Setelah menyalin file ini:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl restart smarthom:*
; =============================================================

; -------------------------------------------------------
; [1] MQTT Listener
; -------------------------------------------------------
[program:laravel-mqtt]
process_name=%(program_name)s
command=/usr/bin/php /www/wwwroot/nh.mdpower.io/smarthom/artisan mqtt:listen
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/www/wwwroot/nh.mdpower.io/smarthom/storage/logs/supervisor-mqtt.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=HOME="/root",USER="root"

; -------------------------------------------------------
; [2] Laravel Scheduler (daemon)
; -------------------------------------------------------
[program:laravel-scheduler]
process_name=%(program_name)s
command=/usr/bin/php /www/wwwroot/nh.mdpower.io/smarthom/artisan schedule:work
autostart=true
autorestart=true
startsecs=3
startretries=5
stopwaitsecs=5
user=root
redirect_stderr=true
stdout_logfile=/www/wwwroot/nh.mdpower.io/smarthom/storage/logs/supervisor-scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
environment=HOME="/root",USER="root"

; -------------------------------------------------------
; [3] Laravel Queue Worker
; -------------------------------------------------------
[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /www/wwwroot/nh.mdpower.io/smarthom/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --queue=default
autostart=true
autorestart=true
numprocs=1
startsecs=5
startretries=5
stopwaitsecs=60
user=root
redirect_stderr=true
stdout_logfile=/www/wwwroot/nh.mdpower.io/smarthom/storage/logs/supervisor-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=HOME="/root",USER="root"

; -------------------------------------------------------
; Group
; -------------------------------------------------------
[group:smarthom]
programs=laravel-mqtt,laravel-scheduler,laravel-queue
priority=999
